# legacy/legacy-printer-app.service.in
[Unit]
Description=Legacy Printer Application

[Service]
ExecStart=legacy-printer-app server -o admin-group=@ADMINGROUP@ -o auth-service=@AUTHSERV@ @LISTENHOST@-o log-file=@LOGFILE@ -o log-level=info -o server-options=@SOPTIONS@
ExecStop=legacy-printer-app shutdown
Type=simple

[Install]
WantedBy=multi-user.target
